# deploy/systemd/escrowhub-ltc-watcher.service
[Unit]
Description=EscrowHub LTC Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/escrowhub
EnvironmentFile=/etc/escrowhub/escrowhub.env
ExecStart=/usr/bin/python3 run_ltc_watcher.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
